# traffic-filter.service
[Unit]
Description=Traffic Filter - eBPF/XDP Based Network Filter
Documentation=https://github.com/your-repo/traffic-filter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# 工作目錄
WorkingDirectory=/opt/traffic-filter

# 執行命令 - 根據需要修改參數
ExecStart=/opt/traffic-filter/traffic-filter \
    -iface eth0 \
    -domains "pornhub.com,www.pornhub.com,xvideos.com,xnxx.com" \
    -dns-mode poison

# 重啟策略
Restart=on-failure
RestartSec=5s

# 資源限制
LimitNOFILE=65536
LimitMEMLOCK=infinity

# 安全設置
NoNewPrivileges=true
PrivateTmp=true

# 日誌
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traffic-filter

[Install]
WantedBy=multi-user.target
